# usb_secure_autoexec/udev/10-autorun.rules
KERNEL=="sd*[!0-9]", ACTION=="add", RUN+="/home/brad/test/test.sh"
